# systemd/rcpd.service
[Unit]
Description=RCP Daemon
After=network.target
Documentation=https://github.com/open-rcp/rcp

[Service]
Type=forking
ExecStart=/usr/local/bin/rcpdaemon start
ExecStop=/usr/local/bin/rcpdaemon stop
Restart=on-failure
RestartSec=5
User=rcp
Group=rcp
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
